# generators/app/templates/src/deb/control/postinst
update-rc.d [[artifactId]] defaults
service [[artifactId]] start || echo "[[artifactId]] could not be started. Try manually with service [[artifactId]] start"

